# pkgbuild/power-profiles-daemon-biglinux.install
post_install() {
systemctl enable --now power-profiles-daemon-biglinux-cpufreq.path
systemctl mask power-profiles-daemon.service
systemctl enable --now power-profiles-daemon-biglinux.service
}

pre_remove() {
systemctl disable --now power-profiles-daemon-biglinux-cpufreq.path
systemctl unmask power-profiles-daemon.service
systemctl disable --now power-profiles-daemon-biglinux.service
} 
